# host/sshjump.service
[Unit]
Description=sshJump
After=network-online.target

[Service]
ExecStart=/opt/sshjump/start_systemd.sh

[Install]
WantedBy=multi-user.target
